# ipk-source/sdxpinn-quecmanager-beta/root/etc/init.d/quecmanager-services
#!/bin/sh /etc/rc.common
START=49
STOP=10
USE_PROCD=1

start_service() {
    # Run clear_logs.sh once without procd management, with a timeout
    echo "Clearing QuecManager Logs"
    timeout 5 /www/cgi-bin/services/clear_logs.sh
    if [ $? -eq 124 ]; then
        echo "clear_logs.sh timed out after 5 seconds"
    else
        echo "Log Cleaner complete"
    fi

    # Start the continuously running services
    echo "Starting AT Command Queue Processor..."
    procd_open_instance
    procd_set_param command /www/cgi-bin/services/at_queue_processor.sh
    procd_set_param respawn
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance
    echo "AT Queue Processor Started"
    echo "Starting Signal Metrics Logger..."
    procd_open_instance
    procd_set_param command /www/cgi-bin/services/log_signal_metrics.sh
    procd_set_param respawn
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance
    echo "Signal Metrics Logger started"

    echo "All QuecManager services Started"
}

stop_service() {
    # procd will handle stopping all instances automatically
    echo "Stopping QuecManager services."
}
