# ipk-source/sdxlemur-factory-packages/root/opt/lib/opkg/info/vsftpd.list
/opt/lib/upgrade/keep.d/vsftpd
/opt/etc/vsftpd/vsftpd.conf
/opt/sbin/vsftpd
/opt/etc/init.d/S49vsftpd
