# simplefirewall/systemd/simplefirewall.service
[Unit]
Description=Simple Firewall Setup
After=network.target

[Service]
Type=oneshot
ExecStart=/usrdata/simplefirewall/simplefirewall.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
